# cron.daily/aquam.sh
#!/bin/bash

cd /root/naver_scraper && /root/Env/aquam/bin/scrapy crawl naverblog